# devicebox.service
[Unit]
Description=DeviceBox - Raspberry Pi Web Interface
After=network.target
Wants=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/devicebox
Environment=PATH=/home/pi/devicebox/venv/bin
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/pi/devicebox/venv/bin/python -u /home/pi/devicebox/app.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Sicherheitseinstellungen
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/pi/devicebox
ReadWritePaths=/tmp

[Install]
WantedBy=multi-user.target
