# devicebox.service
[Unit]
Description=DeviceBox - Raspberry Pi Web Interface
After=network.target
Wants=network.target

[Service]
Type=simple
User=%i
Group=%i
WorkingDirectory=/home/%i/devicebox
Environment=PATH=/home/%i/devicebox/venv/bin
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/%i/devicebox/venv/bin/python -u /home/%i/devicebox/app.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Sicherheitseinstellungen
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/%i/devicebox
ReadWritePaths=/tmp

[Install]
WantedBy=multi-user.target
